# hdd-spindown.service
[Unit]
Description=Automatic Disk Standby

[Service]
Type=simple
ExecStart=/usr/bin/hdd-spindown.sh

[Install]
WantedBy=multi-user.target
